# service-pdf/docker/service/opg-lpa-pdf2/run
#!/bin/sh

rm -rf /tmp/pdf_ramdisk/*
exec su -c "/srv/opg-lpa-pdf2/application/current/vendor/bin/resque >> /var/log/app/queue.log 2>&1" app
